# rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
chmod a+x /etc/init.d/unblockmusic
exit 0
